# deploy/gemhunter-web.service
[Unit]
Description=GemHunter web app on :8080
After=network-online.target

[Service]
Type=simple
User=dkarm
WorkingDirectory=/home/dkarm/Projects/GemHunter
ExecStart=/home/dkarm/Projects/GemHunter/.venv/bin/python -m gemhunter.webapp --db /mnt/ssd/gemhunter/gemhunter.db --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
